# data/systemd/influxdb.service.d/nymea-env.conf
[Service]
Environment="INFLUXDB_HTTP_BIND_ADDRESS=127.0.0.1:8086"
Environment="INFLUXDB_MONITOR_STORE_ENABLED=false"
Environment="INFLUXDB_HTTP_LOG_ENABLED=false"
Environment="INFLUXDB_CONTINUOUS_QUERIES_LOG_ENABLED=false"
Environment="INFLUXDB_CONTINUOUS_QUERIES_QUERY_STATS_ENABLED=false"
Environment="INFLUXDB_DATA_QUERY_LOG_ENABLED=false"
